# semblance/com.csra.health-check.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Local macOS scheduler for the CSRA health check. Runs every 10 minutes
     (StartInterval=600s) and only reruns tests when the tree changed.
     Install with: bash semblance/install-launchd.sh
     Uninstall with: bash semblance/install-launchd.sh --uninstall -->
<plist version="1.0">
<dict>
  <key>Label</key><string>com.csra.health-check</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__ROOT__/semblance/health-check.sh</string>
  </array>
  <key>WorkingDirectory</key><string>__ROOT__</string>
  <key>StartInterval</key><integer>600</integer>
  <key>RunAtLoad</key><true/>
  <key>StandardOutPath</key><string>__ROOT__/semblance/launchd.out.log</string>
  <key>StandardErrorPath</key><string>__ROOT__/semblance/launchd.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key><string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key><string>__HOME__</string>
  </dict>
</dict>
</plist>
